# root/etc/s6-overlay/s6-rc.d/svc-agent/run
#!/usr/bin/with-contenv bash

cd /app/agent-service
exec s6-setuidgid abc node dist/index.js
